# data/99-casper-kbd-backlight.rules
# Give the logged-in user write access to Casper keyboard LED control.
# This avoids the need for root/pkexec on every colour change.
SUBSYSTEM=="leds", KERNEL=="casper::kbd_backlight", RUN+="/bin/chgrp video /sys/class/leds/casper::kbd_backlight/led_control", RUN+="/bin/chmod 0660 /sys/class/leds/casper::kbd_backlight/led_control"
